# build/packaging/rpm/decisionengine.service
[Unit]
Description = Decision EngineService
Wants = network-online.target
After = network-onlilne.target

[Service]
Type = simple
EnvironmentFile = -/etc/sysconfig/decisionengine
ExecStart = /usr/sbin/decision-engine start
ExecStop = /usr/sbin/decision-engine stop
ExecReload = /bin/kill -HUP $MAINPID

[Install]
WantedBy = multi-user.target
