# build/packaging/rpm/decisionengine.service
[Unit]
Description = Decision EngineService
Requires = httpd.service
After = httpd.service

[Service]
Type = simple
ExecStart = /usr/sbin/decision-engine start
ExecStop = /usr/sbin/decision-engine stop
ExecReload = /bin/kill -HUP $MAINPID

[Install]
WantedBy = multi-user.target
